# deploy/systemd/llmlibrarian-chroma@.service
[Unit]
Description=llmLibrarian ChromaDB server (chroma run)
After=network.target
Before=llmlibrarian-mcp.service

[Service]
Type=simple
WorkingDirectory={{INSTALL_DIR}}
ExecStart={{INSTALL_DIR}}/scripts/run_chroma_server.sh
Restart=always
RestartSec=2
Environment=LLMLIBRARIAN_INSTALL_DIR={{INSTALL_DIR}}
Environment=LLMLI_CHROMA_ENV_FILE={{INSTALL_DIR}}/.env.mcp
StandardOutput=append:{{LOG_DIR}}/llmlibrarian-chroma.stdout.log
StandardError=append:{{LOG_DIR}}/llmlibrarian-chroma.stderr.log

[Install]
WantedBy=default.target
