# deploy/systemd/llmlibrarian-mcp@.service
[Unit]
Description=llmLibrarian MCP HTTP service
After=network.target

[Service]
Type=simple
WorkingDirectory={{INSTALL_DIR}}
ExecStart={{INSTALL_DIR}}/scripts/run_mcp_http.sh
Restart=always
RestartSec=2
Environment=LLMLI_MCP_ENV_FILE={{INSTALL_DIR}}/.env.mcp

[Install]
WantedBy=default.target
