# deploy/launchd/com.tjm4.llmlibrarian-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tjm4.llmlibrarian-mcp</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/tjm4/llmLibrarian/scripts/run_mcp_http.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/tjm4/llmLibrarian</string>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/tjm4/Library/Logs/llmlibrarian-mcp.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/tjm4/Library/Logs/llmlibrarian-mcp.stderr.log</string>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
